# deploy/notes-app.service.example
# systemd service unit for the Notes App (Next.js)
# ─────────────────────────────────────────────────────────────────
# 1. Copy to /etc/systemd/system/notes-app.service
# 2. Edit User, WorkingDirectory, and EnvironmentFile to match your setup.
# 3. sudo systemctl daemon-reload
# 4. sudo systemctl enable --now notes-app
# ─────────────────────────────────────────────────────────────────

[Unit]
Description=Notes App — Next.js Service
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/clawapp
EnvironmentFile=/path/to/clawapp/.env.local
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
